# piclock-customisations/06-display-change/files/99-display-change.rules
KERNEL=="card*", SUBSYSTEM=="drm", ACTION=="change", RUN+="/usr/local/sbin/display-change-handler"
